# scripts/com.curbscout.pipeline.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.curbscout.pipeline</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/env</string>
        <string>bash</string>
        <string>-c</string>
        <string>cd ~/antigravity/CurbScout/pipeline &amp;&amp; source .venv/bin/activate &amp;&amp; python -m curbscout.main</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>/tmp/curbscout-pipeline.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/curbscout-pipeline.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>GCP_HUB_URL</key>
        <string>https://curbscout-hub-HASH-uc.a.run.app</string>
        <key>WORKER_API_KEY</key>
        <string>REPLACE_WITH_REAL_KEY</string>
        <key>GCS_BUCKET</key>
        <string>curbscout-artifacts</string>
    </dict>
</dict>
</plist>
